# nas-web.service
[Unit]
Description=NAS web desktop & setup wizard (nas-web.py)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
Environment=NAS_WEB_PORT=80
Environment=SUDO_USER=oleg
WorkingDirectory=/home/oleg
ExecStart=/usr/bin/python3 /home/oleg/nas-web.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
